# ops/deliverty-api.service.example
[Unit]
Description=Deliverty API
After=network.target
# Note: Production uses Supabase Postgres, not local PostgreSQL

[Service]
Type=simple
# Set DATABASE_URL to your Supabase Postgres connection string
# Get it from: Supabase Dashboard → Project Settings → Database → Connection string
# Example: postgresql://postgres:[YOUR-PASSWORD]@[PROJECT].supabase.co:5432/postgres?sslmode=require
Environment=DATABASE_URL=postgresql://postgres:[PASSWORD]@[PROJECT].supabase.co:5432/postgres?sslmode=require
Environment=HTTP_ADDR=:8080
Environment=TG_BOT_TOKEN=your-bot-token-here
Environment=TG_BOT_NAME=deliverty_bot
Environment=TG_DEEPLINK_SECRET=your-secret-here
ExecStart=/srv/deliverty-api/api
Restart=always
RestartSec=5
User=www-data
Group=www-data
WorkingDirectory=/srv/deliverty-api

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deliverty-api

[Install]
WantedBy=multi-user.target
